# crontab/syncdirs.sh
#!/bin/bash

rsync -av --delete $SRC_DIR_PATH $DST_DIR_PATH