# crontab/syncdirs.sh
#!/bin/bash

rsync -aq --delete $SRC_DIR_PATH $DST_DIR_PATH